# crontab/stats.crontab
SHELL=/bin/bash
0 10 1 * * costanza /bin/bash -c '/home/costanza/costanza report --month $(date --date="$(date +%Y-%m-15) - 1 month" "+%Y%m")' >> "/var/log/costanza/report-$(date +%Y-%m).log"
0 11 1 * * costanza /bin/bash -c '/home/costanza/costanza report remove --month $(date --date="$(date +%Y-%m-15) - 1 month" "+%Y%m")' >> "/var/log/costanza/remove-report-$(date +%Y-%m).log"
